# deploy/systemd/nbdpsy-mcp.service
[Unit]
Description=nbdpsy-mcp - 小红书运营能力 MCP 后台
After=network.target nbdpsy-xvfb.service
Requires=nbdpsy-xvfb.service

[Service]
Type=simple
User=roots
WorkingDirectory=/home/roots/nbdpsy-mcp
Environment=DISPLAY=:99
ExecStartPre=/home/roots/nbdpsy-mcp/.venv/bin/alembic upgrade head
ExecStartPre=/bin/bash /home/roots/nbdpsy-mcp/scripts/pack_extension.sh
ExecStart=/home/roots/nbdpsy-mcp/.venv/bin/uvicorn app.server:create_app --factory --host 0.0.0.0 --port 8848
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
